# dist/openrc/xvfb.initd
#!/sbin/runscript
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

start() {
    ebegin "Starting virtual X frame buffer: Xvfb"
    start-stop-daemon --start --pidfile ${PIDFILE} --make-pidfile --background --exec ${XVFB} -- ${XVFBARGS}
    eend $?
}

stop() {
    ebegin "Stopping virtual X frame buffer: Xvfb"
    start-stop-daemon --stop --pidfile ${PIDFILE}
    eend $?
}
