# dist/systemd/url2img.service
[Unit]
Description=Capture screenshots of websites

[Service]
User=daemon
Group=daemon
Environment=DISPLAY=:1
PIDFile=/run/url2img.pid
ExecStart=/usr/local/bin/url2img -cache-dir /var/cache/url2img
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
#MemoryLimit=1536M

[Install]
WantedBy=multi-user.target
