# dist/systemd/url2img.service
[Unit]
Description=Capture screenshots of websites

[Service]
User=daemon
Group=daemon
PIDFile=/run/url2img.pid
ExecStartPre=/usr/bin/mkdir -p /var/cache/url2img
ExecStart=/usr/local/bin/url2img -cache-dir /var/cache/url2img
ExecReload=/usr/bin/kill -HUP $MAINPID
Restart=on-abnormal
#MemoryLimit=1536M

[Install]
WantedBy=multi-user.target
